# crontab
0 0 * * * cd /app && echo "$(date) - Running session_cleanup" && flask session_cleanup 2>&1
0 0 * * * cd /app && echo "$(date) - Running refresh-puzzles" && flask refresh-puzzles 2>&1

